# usr/src/oss/postfix/files/postfix.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">

<service_bundle type='manifest' name='postfix'>

<service
       name='network/postfix'
       type='service'
       version='1'>

        <create_default_instance enabled='false'/>
        <single_instance/>

        <dependency name='config-file'
                    grouping='require_all'
                    restart_on='none'
                    type='path'>
                <service_fmri value='file:///etc/postfix/main.cf'/>
        </dependency>

       <dependency name='loopback'
                   grouping='require_all'
                   restart_on='error'
                   type='service'>
               <service_fmri value='svc:/network/loopback:default'/>
       </dependency>

       <dependency name='physical'
                   grouping='require_all'
                   restart_on='error'
                   type='service'>
               <service_fmri value='svc:/network/physical:default'/>
       </dependency>

       <dependency name='fs-local'
                   grouping='require_all'
                   restart_on='none'
                   type='service'>
               <service_fmri value='svc:/system/filesystem/local'/>
       </dependency>

       <exec_method
               type='method'
               name='start'
               exec='/usr/sbin/postfix start'
               timeout_seconds='60' />

       <exec_method
               type='method'
               name='stop'
               exec='/usr/sbin/postfix stop'
               timeout_seconds='60' />

       <exec_method
               type='method'
               name='refresh'
               exec='/usr/sbin/postfix reload'
               timeout_seconds='60' />

       <stability value='Unstable' />

       <template>
               <common_name>
                       <loctext xml:lang='C'>Postfix SMTP Server</loctext>
               </common_name>
               <documentation>
                       <manpage title='postfix' section='1' />
                       <doc_link name='postfix.org'
                               uri='http://www.postfix.org/documentation.html' />
               </documentation>
       </template>
</service>

</service_bundle>
